# systemd/jamie-transcribe.service
[Unit]
Description=Jamie Transcription Service
After=network.target postgresql.service

[Service]
ExecStart=/usr/local/bin/jamie transcribe
Restart=always
User=jamie
Group=jamie
Environment=DATABASE_URL=postgres://jamie:${POSTGRES_PASSWORD}@localhost:5432/jamie
Environment=DISCORD_TOKEN=${DISCORD_TOKEN}
Environment=GEMINI_API_KEY=${GEMINI_API_KEY}
Environment=SPEECHMATICS_API_KEY=${SPEECHMATICS_API_KEY}

[Install]
WantedBy=multi-user.target
